# templates/apprise_healthcheck.cron
#birdnet
0 1 * * 1 $USER /usr/local/bin/apprise_healthcheck.sh >/dev/null 2>&1
