# templates/apprise_healthcheck.cron
#birdnet
*/15 * * * * $USER /usr/local/bin/apprise_healthcheck.sh >/dev/null 2>&1
