# templates/automatic_update.cron
#birdnet
7 0 * * * $USER /usr/local/bin/automatic_update.sh >/dev/null 2>&1
